# sbomguard.service
[Unit]
Description=SBOMguard — SBOM vulnerability monitor
After=network.target

[Service]
Type=simple
User=soc
WorkingDirectory=$HOME/claude/SBOMguard
EnvironmentFile=$HOME/claude/SBOMguard/.env
ExecStart=/usr/bin/python3 $HOME/claude/SBOMguard/app.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
