# soc-sbom.service
[Unit]
Description=SBOMguard — SBOM vulnerability monitor
After=network.target

[Service]
Type=simple
User=soc
WorkingDirectory=$HOME/soc-sbom
EnvironmentFile=$HOME/soc-sbom/.env
ExecStart=/usr/bin/python3 $HOME/soc-sbom/app.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
